# ttsalert.service
[Unit]
Description=TTS Alert Service - Voice Alert System
Documentation=https://github.com/jh163888/ttsalert
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ttsalert
Group=ttsalert

ExecStart=/usr/local/bin/ttsalert
Restart=on-failure
RestartSec=5s

Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/var/lib/ttsalert"

WorkingDirectory=/var/lib/ttsalert

ReadWritePaths=/var/lib/ttsalert/audio
ReadWritePaths=/etc/ttsalert

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ttsalert

[Install]
WantedBy=multi-user.target
